# deploy/pi-dashboard-kiosk.service
[Unit]
Description=Pi Dashboard Kiosk Browser
After=graphical.target display-manager.service network-online.target pi-dashboard.service
Wants=graphical.target network-online.target pi-dashboard.service
StartLimitIntervalSec=0

[Service]
Type=simple
User={{APP_USER}}
WorkingDirectory={{APP_DIR}}
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/{{APP_USER}}/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/%U
ExecStart={{APP_DIR}}/deploy/kiosk.sh http://localhost:8080
Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target
